# config/crontab
0 2 * * * python /app/batch_process.py >> /app/logs/batch.log 2>&1